# infrastructure/systemd/osint-platform-local.service
[Unit]
Description=OSINT Platform Local Production Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/rfileto/projects/osint-platform
ExecStart=/usr/bin/docker compose --env-file /home/rfileto/projects/osint-platform/.env -f /home/rfileto/projects/osint-platform/docker-compose.yml -f /home/rfileto/projects/osint-platform/compose.prod.yml up -d
ExecStop=/usr/bin/docker compose --env-file /home/rfileto/projects/osint-platform/.env -f /home/rfileto/projects/osint-platform/docker-compose.yml -f /home/rfileto/projects/osint-platform/compose.prod.yml down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
